# peervpn.initd
#!/sbin/openrc-run

# Sample init.d file for alpine linux.

name=peervpn
command="/usr/sbin/$name"
command_args="/etc/conf.d/$name"
command_background="yes"

# peervpn does not daemonize itself.
start_stop_daemon_args="--background pidfile=/run/peervpn.pid --make-pidfile --remove-pidfile"
pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}
